# monitor_agent/supervisord.conf
[supervisord]
http_port=/var/tmp/supervisor.sock ; (default is to run a UNIX domain socket server)
;http_port=127.0.0.1:9001  ; (alternately, ip_address:port specifies AF_INET)
;sockchmod=0700              ; AF_UNIX socketmode (AF_INET ignore, default 0700)
;sockchown=nobody.nogroup     ; AF_UNIX socket uid.gid owner (AF_INET ignores)
;umask=022                   ; (process file creation umask;default 022)
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;http_username=user          ; (default is no username (open system))
;http_password=123           ; (default is no password (open system))
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;user=chrism                 ; (default is current user, required if root)
;directory=/tmp              ; (default is not to cd during start)
;environment=KEY=value       ; (key value pairs to add to environment)

[supervisorctl]
serverurl=unix:///var/tmp/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:ksmonitor]
command=python /opt/LvsMonitorAgent/ksmonitor.py           
autostart=true            
autorestart=true          
startsecs=1              
startretries=3            
exitcodes=0,2             
stopsignal=QUIT           
stopwaitsecs=10           
user=chrism               
log_stdout=true           
log_stderr=true           
logfile=/opt/LvsMonitorAgent/data/ksmonitor.log  
logfile_maxbytes=1MB      
logfile_backups=10        

[program:run.py]
command=python /opt/LvsMonitorAgent/run.py           
autostart=true            
autorestart=true          
startsecs=10              
startretries=3            
exitcodes=0,2             
stopsignal=QUIT           
stopwaitsecs=10           
user=chrism               
log_stdout=true           
log_stderr=true           
logfile=/opt/LvsMonitorAgent/data/run.log  
logfile_maxbytes=1MB      
logfile_backups=10 

[program:keepalived_log_analyze.py]
command=python /opt/LvsMonitorAgent/keepalived_log_analyze.py         
autostart=true            
autorestart=true          
startsecs=10            
startretries=3           
exitcodes=0,2             
stopsignal=QUIT           
stopwaitsecs=10
user=chrism               
log_stdout=true           
log_stderr=true           
logfile=/opt/LvsMonitorAgent/data/keepalived_log_analyze.log  
logfile_maxbytes=1MB      
logfile_backups=10 